# config/s6-overlay/s6-rc.d/openvpn-exporter/run
#!/usr/bin/env bash

exec /usr/src/app/openvpn-exporter/bin/start.sh
