# config/s6-overlay/s6-rc.d/open-balena-vpn/run
#!/usr/bin/env bash
# shellcheck disable=SC1091

set -euo pipefail

# Redirect all future stdout/stderr to s6-log
exec > >(exec s6-log -b p"vpn[$$]:" 1 || true) 2>&1

# Change to working directory
cd /usr/src/app || exit 1

# Load environment variables for this service
source /etc/s6-overlay/scripts/functions.sh
[[ -f "config/env" ]] && load_env_file "config/env"

# ShellCheck is currently bad at figuring out functions that are invoked via trap.
# In such cases, please ignore the message with a directive.
# shellcheck disable=SC2329,SC2317
signal_handler() {
    
    # Run the backend VPN workers draining logic
    echo "Draining backend VPN workers..."
    /etc/s6-overlay/scripts/stop.sh || true # ignore errors from this script
    
    # Forward SIGTERM to main process
    if [[ -n ${main_pid} ]] && kill -0 "${main_pid}" 2>/dev/null; then
        echo "Forwarding SIGTERM to main process (PID: ${main_pid})"
        kill -TERM "${main_pid}" 2>/dev/null || true
    fi
    
    # Wait for main process to finish its graceful shutdown
    wait "${main_pid}"
    exit $?
}

# Trap signals before starting the service (killmode=mixed)
trap 'signal_handler' SIGTERM SIGINT

# Start service in background (no exec)
/etc/s6-overlay/scripts/start.sh &
main_pid=$!

# Wait for the process
wait "${main_pid}"
