# config/services/confd.service.d/confd.conf
[Service]
ExecStartPre=-/bin/bash -ec 'systemctl set-environment NPROC=$(nproc)'
