# config/s6-overlay/s6-rc.d/extra-contenv/up
/etc/s6-overlay/scripts/extra-contenv-up
