# config/s6-overlay/s6-rc.d/haproxy/run
#!/usr/bin/env bash
# shellcheck disable=SC1091

# HAProxy service for s6-overlay

set -euo pipefail

# Redirect all future stdout/stderr to s6-log
exec > >(exec s6-log -b p"haproxy[$$]:" 1 || true) 2>&1

# Change to working directory
cd /usr/src/app || exit 1

# Load environment variables for this service
source /etc/s6-overlay/scripts/functions.sh
[[ -f "config/env" ]] && load_env_file "config/env"

# Load environment files (optional configs)
[[ -f /etc/default/haproxy ]] && source /etc/default/haproxy
[[ -f /etc/sysconfig/haproxy ]] && source /etc/sysconfig/haproxy

# Set default environment variables
CONFIG="/etc/haproxy/haproxy.cfg"
PIDFILE="/run/haproxy.pid"
EXTRAOPTS=("-S" "/run/haproxy-master.sock")

# Validate configuration before starting
if ! /usr/sbin/haproxy -f "${CONFIG}" -c -q "${EXTRAOPTS[@]}"; then
    exit 1
fi

# -W: master-worker mode
# -db: disable background/daemon mode (run in foreground for s6)
# -f: config file
# -p: pidfile
# Run haproxy in foreground with exec to pass signals through
exec /usr/sbin/haproxy -W -db -f "${CONFIG}" -p "${PIDFILE}" "${EXTRAOPTS[@]}"
