# config/services/haproxy.service.d/haproxy.conf
[Unit]
Requires=confd.service
After=confd.service

[Service]
EnvironmentFile=/usr/src/app/config/env

[Install]
WantedBy=
WantedBy=basic.target
