# config/services/haproxy.service.d/haproxy.conf
[Unit]
Requires=confd.service
After=confd.service

[Install]
WantedBy=
WantedBy=basic.target
