# config/services/openvpn@.service.d/namespace.conf
[Service]
PrivateTmp=no
ProtectSystem=no
